# twit-dashboard.service
[Unit]
Description=TWiT E-Ink Dashboard renderer and delivery

[Service]
Type=oneshot
ExecStart=/home/leo/Projects/twit-eink-dashboard/.venv/bin/python3 /home/leo/Projects/twit-eink-dashboard/dashboard.py
WorkingDirectory=/home/leo/Projects/twit-eink-dashboard
